# appliance/etc/cron.d/update-certificate
0 0 1 * * /scripts/update-certificate.sh
